# librecad/supervisord.conf
[program:librecad]
command=/opt/wine-stable/bin/wine 'C:\Program Files\LibreCAD\LibreCAD.exe'
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true
